# debian/preinst
#! /bin/sh 
# preinst script for #PACKAGE#
#
# see: dh_installdeb(1)
#
# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
#
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

set -e

set_config_value()
{
    file=$1
    option=$2
    value=$3

    text=$(sed "s/\($option[[:space:]]*=[[:space:]]*\).*/\1$value/" "$file") &&
        echo "$text" > "$file" &&
	echo "" >> "$file"
}



case "$1" in
    install|upgrade)
	### If we are upgrading from v0.7, the "UGID=" setting in
        ### /etc/init.d/greylist would have been overwritten by DebConf.
        ### Now, however, this is a conffile.  Set this value back to its
        ### default to restore the script to its pristine state, so that
        ### dpkg won't bother the user about the modified conffile.

	if dpkg --compare-versions "$2" eq "0.7"
	then
	    set_config_value /etc/init.d/greylist UGID mail:mail
	fi

	### In version 0.8.8 the /etc/init.d/ file is called greylist.
        ### Newer packages call this file greylistd, so manually remove
        ### the old one.
	if dpkg --compare-versions "$2" eq "0.8.8"
	then
	    echo "I: update-rc.d -f greylist remove"
	    update-rc.d -f greylist remove
	fi
    ;;

    abort-upgrade)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


